# gunicorn/supervisor.conf
[program:{{ name }}]
command={{ venv_path }}/bin/gunicorn {{ app }} -b unix:{{ socket }} -w {{ num_workers }}
directory={{ program_dir }}
user={{ user }}
autostart=true
autorestart=true
