# saltbot/supervisor.conf
[program:saltbot]
command=exec /home/saltbot/venv/bin/saltbot /home/saltbot/saltbot.yml
directory=/home/saltbot
user=saltbot
autostart=true
autorestart=true
stopasgroup=true
stopsignal=KILL
