# monitoring/sentry/sentry-workers-supervisord.conf
[program:sentry-workers]
directory=/home/sentry
command=/home/sentry/venv/bin/sentry celery worker -B
user=sentry
autostart=true
autorestart=true
