# habhub/supervisor.conf
[program:dfvc]
command=/home/dfvc/venv/bin/gunicorn app:app -b localhost:8001 -w 2
directory=/home/dfvc/dl-fldigi/update_server
user=dfvc
autostart=true
autorestart=true
