# tawhiri/tawhiri-api.supervisor.conf
[program:tawhiri-api]
directory=/srv/tawhiri
user=tawhiri
autostart=true
autorestart=true
command=/srv/tawhiri/bin/gunicorn --config /srv/tawhiri/gunicorn_cfg.py tawhiri.api:app
